# common/systemd/docker-redeploy-edge.service
[Unit]
Description=Redeploy Docker stack on edge after git sync
After=git-sync.service docker.service
Requires=git-sync.service docker.service

[Service]
Type=oneshot
User=root
Group=docker
WorkingDirectory=/opt/home-lab/edge
Environment=COMPOSE_PROJECT_NAME=edge
EnvironmentFile=-/opt/home-lab/edge/.env
ExecStart=/usr/bin/docker compose up -d
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
