# common/systemd/docker-redeploy-edge.service
[Unit]
Description=Redeploy Docker stack on edge server after git sync
After=git-sync.service
Requires=git-sync.service

[Service]
Type=oneshot
WorkingDirectory=/opt/home-lab/edge
ExecStart=/usr/local/bin/docker compose up -d
User=root
Group=docker

[Install]
WantedBy=multi-user.target
