# common/systemd/docker-redeploy-homelab.service
[Unit]
Description=Redeploy Docker stack on homelab after git sync
After=git-sync.service
Requires=git-sync.service

[Service]
Type=oneshot
WorkingDirectory=/opt/home-lab/homelab
ExecStart=/bin/bash -c 'find . -name "docker-compose.yml" -execdir /usr/local/bin/docker compose up -d \;'
User=root
Group=docker

[Install]
WantedBy=multi-user.target
