# common/systemd/git-sync.service
[Unit]
Description=Git sync service for home-lab repository
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/home-lab
ExecStart=/opt/home-lab/common/scripts/git-sync.sh
StandardOutput=journal
StandardError=journal