# edge/systemd/caddy-reload.service
[Unit]
Description=Regenerate Caddyfile from Docker labels
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/home-lab/edge
ExecStart=/usr/bin/bash generate-caddyfile.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
